# meta-rcar-gen3/recipes-rswitch/rswitch1-startup/files/rswitch-init.target
[Unit]
Description=Renesas R-Switch Initialization
Requires=network.target
After=network.target
